# crontab
*/5 * * * * /home/ubuntu/nms-python-telegram/run.sh
